# deploy/ansible/roles/services/templates/upstart.conf
#!upstart
#
# Copy me to /etc/init/
# This installs a daemon as a system level call and ensures the process is
# consistently restarted on error.
# Manual start, stop, and restart respected.
#

description "Daemon for nodejs black-phoenix server"

start on startup
stop on shutdown
respawn

setuid {{app_user}}
setgid {{app_user}}

env NODE_ENV={{node_env}}

script
  export DEBUG=*

  /usr/bin/node {{app_base_path}} >> {{log_file}} 2>&1
end script
